# dev/quadlets/mars-seed.container
[Unit]
Description=MARS local OSM seed
Requires=mars-postgis.service
After=mars-postgis.service

[Container]
Image=localhost/mars-seed:dev
ContainerName=mars-seed
Network=mars.network
Volume=osm-cache.volume:/cache
Environment=PG_DSN=postgres://mars:mars@mars-postgis/mars OSM_BBOX=55.6,9.5,56.0,10.7

[Service]
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=default.target
