# dev/quadlets/mars-postgis.container
[Unit]
Description=MARS local PostGIS

[Container]
Image=docker.io/postgis/postgis:18-3.6
ContainerName=mars-postgis
Network=mars.network
Volume=postgis-data.volume:/var/lib/postgresql/data
Volume=%h/.local/share/mars-dev/postgis/init.sql:/docker-entrypoint-initdb.d/00-init.sql:ro,Z
Environment=POSTGRES_DB=mars POSTGRES_USER=mars POSTGRES_PASSWORD=mars
Exec=postgres -c wal_level=logical -c max_replication_slots=8 -c max_wal_senders=8
HealthCmd=pg_isready -U mars -d mars
HealthInterval=2s
HealthRetries=30

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
